# etc/systemd/system/spotbot_script.timer
[Unit]
Description=Run spotipi_script.service every 2s

[Timer]
AccuracySec=1s
OnBootSec=5sec
OnUnitActiveSec=2s